# patches/owner-console-v3.0.0/drivelab-owner-ops.path
[Unit]
Description=Watch for allowlisted DriveLab Owner Console operations

[Path]
PathExistsGlob=/var/lib/drivelab-license/ops/requests/*.request
Unit=drivelab-owner-ops.service

[Install]
WantedBy=multi-user.target
